# llm_server.service
[Unit]
Description=LLM Server (Uvicorn)
After=network.target

[Service]
# Path to your project directory
WorkingDirectory=/home/eb3-brayan/eb3_llm_server

# Activate your virtual environment
ExecStart=/home/eb3-brayan/eb3_llm_server/venv/bin/uvicorn llm_server:app --host 0.0.0.0 --port 8000

# Restart automatically if it crashes
Restart=always
RestartSec=5

# Run as your user
User=eb3-brayan
Group=eb3-brayan

# Logging
StandardOutput=append:/var/log/llm_server.log
StandardError=append:/var/log/llm_server_error.log

[Install]
WantedBy=multi-user.target
